# nginx.conf
[supervisord]
nodaemon=true

[program:nginx]
directory=/etc/nginx
command=nginx -g
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true