# supervisor_simulator.conf
[supervisord]
nodaemon=true

[program:simulator]
directory=/app
command=dotnet sumulator/Belane.IOT.Simulator.dll
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true